# library/user/pager2pager/p2p_pager/init_p2p_pager
#!/bin/sh /etc/rc.common

USE_PROCD=1
START=99
STOP=10

start_service(){
    procd_open_instance
    procd_set_param command /bin/bash "/usr/bin/start_p2p_pager.sh"
    procd_set_param respawn
    procd_close_instance
}

stop_service(){
    procd_killall "python3 /usr/bin/p2p_pager"
}
