# library/user/pager2pager/p2p_pager/init_p2p_pager
#!/bin/sh /etc/rc.common

START=99

start_service(){
    procd_open_instance
    procd_set_param command python3 /usr/bin/p2p_pager
    procd_set_param respawn
    procd_close_instance
}

stop_service(){
    procd_killall "python3 /usr/bin/p2p_pager"
}